# shell/S51arkkvmapp
# description: ArkKVMApp Startup Script

### BEGIN INIT INFO
# Provides:          ArkKVMApp
# Required-Start:    $local_fs $network $remote_fs
# Required-Stop:     $local_fs $network $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start ArkKVMApp at boot time
# Description:       Enable service provided by ArkKVMApp.
### END INIT INFO
APP_PATH=/oem/usr/bin/arkkvm_app
# APP_PATH=/userdata/arkkvm/bin/arkkvm_app
APP_NAME=ArkKVMApp
PID_FILE=/run/$APP_NAME.pid
FT_IMG_PATH="/userdata/arkkvm/imgs/ft_disk.img"

[ -f /etc/profile.d/RkEnv.sh ] && source /etc/profile.d/RkEnv.sh

case "$1" in
    start)
        echo "Starting $APP_NAME..."
        if [ -f $PID_FILE ]; then
            if ps -p $(cat $PID_FILE) > /dev/null 2>&1; then
                echo "$APP_NAME is already running (PID: $(cat $PID_FILE))"
                exit 0
            else
                echo "Removing stale PID file..."
                rm -f $PID_FILE
            fi
        fi

        mkdir -p "$(dirname "$FT_IMG_PATH")"
        if [ ! -f "$FT_IMG_PATH" ]; then
            echo "Creating FAT32 image file: $FT_IMG_PATH"
                dd if=/dev/zero of="$FT_IMG_PATH" bs=1M count=0 seek=16384
                mkfs.vfat -F 32 "$FT_IMG_PATH"
            echo "FAT32 image created successfully."
        else
            echo "Image file already exists: $IMG_PATH"
        fi

        mkdir -p /userdata/arkkvm/log

        nohup sh -c "sleep 5 && $APP_PATH" >> /userdata/arkkvm/log/arkkvm_app.log 2>&1 &
        echo $! > $PID_FILE
        echo "$APP_NAME started with PID $!. /userdata/arkkvm/log/arkkvm_app.log"
        ;;
    stop)
        echo "Stopping $APP_NAME..."
        if [ -f $PID_FILE ]; then
            PID=$(cat $PID_FILE)
            if ps -p $PID > /dev/null 2>&1; then
                kill $PID
                echo "$APP_NAME stopped (PID: $PID)"
            else
                echo "No running process found for PID $PID"
            fi
            rm -f $PID_FILE
        else
            echo "PID file not found: $PID_FILE"
        fi
        ;;
    restart)
        $0 stop
        sleep 2
        $0 start
        ;;
    status)
        if [ -f $PID_FILE ]; then
            PID=$(cat $PID_FILE)
            if ps -p $PID > /dev/null 2>&1; then
                echo "$APP_NAME is running (PID: $PID)"
            else
                echo "$APP_NAME PID file exists but process not running"
            fi
        else
            echo "$APP_NAME is not running"
        fi
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac

exit 0
